# infra/systemd/skysealens-api.service
[Unit]
Description=SkySeaLens API
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/skysealens/current
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=512
EnvironmentFile=/srv/skysealens/shared/.env
ExecStart=/usr/bin/node services/api/src/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
